# debian/adaptation-vayu-configs/etc/systemd/system/phosh.service.d/90-wait.conf
# FIXME
[Service]
ExecStartPre=/usr/bin/sleep 5
